# packages/srcmap-reverser/systemd-units/srcmap-reverse.socket
[Unit]
Description=Source Map Reverser Socket

[Socket]
ListenStream=/var/run/srcmap-reverse.sock

[Install]
WantedBy=sockets.target
